# redis/etc/service/redis/run
#!/bin/sh
cd /redis-data
exec redis-server
exec redis-sentinel /etc/sentinel.conf
