# redis-sentinel/etc/service/redis-sentinel/run
#!/bin/sh
exec redis-sentinel /sentinel.conf
